# stark.service
[Unit]
Description=Stark Print Server
After=network.target

[Service]
User=stark
Group=stark
WorkingDirectory=/home/stark
ExecStart=/home/stark/venv/bin/python /home/stark/server.py
Restart=always
Environment="PATH=/home/stark/venv/bin"

[Install]
WantedBy=multi-user.target

2️⃣ Reload and start

sudo systemctl daemon-reload
sudo systemctl enable stark.service
sudo systemctl start stark.service